# smf/manifests/medusa.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright 2016 Joyent, Inc.
-->

<service_bundle type="manifest" name="smartdc-medusa">
    <service
        name="smartdc/application/medusa"
        type="service"
        version="1" >

        <dependency
            name="network"
            grouping="require_all"
            restart_on="error"
            type="service" >
            <service_fmri
                value="svc:/network/physical" />
        </dependency>

        <dependency
            name="filesystem"
            grouping="require_all"
            restart_on="error"
            type="service" >
            <service_fmri
                value="svc:/system/filesystem/local" />
        </dependency>

        <dependency
            name="config-agent"
            grouping="optional_all"
            restart_on="none"
            type="service">
            <service_fmri
                value="svc:/smartdc/application/config-agent" />
        </dependency>

        <exec_method
            type="method"
            name="start"
            exec="@@NODE@@ --abort-on-uncaught-exception @@PREFIX@@/lib/reflector.js &amp;"
            timeout_seconds="30" >

            <method_context>

                <method_credential
                    user="medusa"
                    group="medusa"
                    privileges="basic,net_privaddr,file_dac_read,file_dac_search" />

                <method_environment>
                    <envvar
                        name="CONFIG_FILE"
                        value="@@PREFIX@@/etc/config.json" />
                </method_environment>

            </method_context>

        </exec_method>

        <exec_method
            type="method"
            name="stop"
            exec=":kill"
            timeout_seconds="30" />

        <instance
            name="default"
            enabled="true" />

        <stability
            value='Unstable' />

        <template>
            <common_name>
                <loctext xml:lang="C">Medusa</loctext>
            </common_name>
        </template>

    </service>
</service_bundle>
<!-- vim: set sw=4 ts=4 sts=4 et: -->
